# home/@images/laravel/crontab
* * * * * cd /var/www/html && /usr/local/bin/php artisan schedule:run --no-ansi > /var/www/html/storage/logs/cron.log 2>&1
0 0 * * * find /tmp -type f -atime +10 -delete
